# openvpn.ini
[program:openvpn]
command = /usr/bin/openvpn --cd /etc/openvpn --config vpn.ovpn --mute-replay-warnings
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
