# transmission.ini
[program:transmission]
user = transmission
group = transmission
command = /usr/bin/transmission-daemon -f -g /var/lib/transmission --no-portmap --rpc-bind-address 0.0.0.0 --bind-address-ipv4 127.0.0.1
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
